# temidara.service
[Unit]
Description=Temidara Dev Rust App
After=network.target

[Service]
# User configuration - change if your user is not 'pi'
User=pi
Group=pi

# Path configuration - change if your repo is elsewhere
WorkingDirectory=/home/pi/temidaradev-rust

# Environment setup
# we add .cargo/bin to PATH so cargo and cargo-leptos can be found
Environment=PATH=/home/pi/.cargo/bin:/usr/local/bin:/usr/bin:/bin
# Bind to 0.0.0.0 so it is accessible outside the Pi
Environment=LEPTOS_SITE_ADDR=0.0.0.0:3004

# Run the command from Makefile: npm run build; cargo leptos watch
# We wrap in bash to execute the compound command
ExecStart=/bin/bash -c 'npm run build && /home/pi/.cargo/bin/cargo leptos watch'

Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
